# apps/infrastructure/ec2/glassact-backup.timer
[Unit]
Description=Run GlassAct Postgres backup nightly

[Timer]
OnCalendar=*-*-* 03:00:00
RandomizedDelaySec=15m
Persistent=true

[Install]
WantedBy=timers.target
