# apps/infrastructure/ec2/glassact-backup.service
[Unit]
Description=GlassAct nightly Postgres backup
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
EnvironmentFile=/opt/glassact/backup.env
ExecStart=/opt/glassact/backup.sh
